# docker/etc/supervisor/conf.d/php-ffmpeg.conf
[program:php-ffmpeg]
command = /usr/bin/php /var/www/index.php
user = www-data
autorestart = true
stderr_logfile = NONE
stdout_logfile = NONE
